# res/systemd/arch-update-manager-check.timer
[Unit]
Description=Periodic check for pending Arch updates (arch-update-manager)

[Timer]
OnCalendar=daily
Persistent=true
Unit=arch-update-manager-check.service

[Install]
WantedBy=timers.target
