# res/systemd/daim-tray.service
[Unit]
Description=Arch Install Manager tray applet
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/daim-tray
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
